# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
sleep 30;
ip link set down phy1-ap0;
ip link set down phy2-ap0;
sleep 5;
ip link set up phy1-ap0;
ip link set up phy2-ap0;
exit 0
