# files/etc/uci-defaults/uci-settings
uci set system.@system[0].timezone=CST-8
uci set system.@system[0].zonename=Asia/Shanghai
uci set system.@system[0].hostname='link-to-the-future'
uci commit system

uci set network.globals.ula_prefix='fd0f:7a7c:a8c9::/48'
uci del network.wan6
uci set firewall.@zone[1].network='wan'
uci set network.wan.peerdns='0'
uci add_list network.wan.dns='::1'
uci set network.wan.proto='pppoe'
uci set network.wan.ipv6='auto'

# Enable wifi by default
uci set wireless.radio0.country='PA'
uci set wireless.radio1.country='PA'
uci set wireless.radio0.disabled=0
uci set wireless.radio1.disabled=0
uci set wireless.radio1.htmode='HE40'
uci set wireless.radio0.htmode='HE160'
uci set wireless.default_radio0.ssid='Link to the FUTURE'
uci set wireless.default_radio1.ssid='Xiaomi_391F'
uci set wireless.default_radio1.encryption='psk-mixed'
uci set wireless.default_radio1.key='1qaz@WSX'
uci set wireless.default_radio0.encryption='sae'
uci set wireless.default_radio0.key='1qaz@WSX'
uci set wireless.radio1.txpower='17'
uci set wireless.radio0.txpower='27'
uci set wireless.radio0.channel='44'
uci set wireless.radio0.cell_density='0'
uci set wireless.radio1.cell_density='0'
uci set wireless.default_radio0.dtim_period='3'
uci set wireless.default_radio0.max_inactivity='86400'
uci set wireless.default_radio0.disassoc_low_ack='0'
uci set wireless.default_radio1.dtim_period='3'
uci set wireless.default_radio1.max_inactivity='86400'
uci set wireless.default_radio1.disassoc_low_ack='0'
# Force 40Mhz on 2.4Ghz
uci set wireless.radio1.noscan=1
uci set wireless.radio1.channel='8'

uci set dhcp.odhcpd.leasetrigger='/usr/lib/unbound/odhcpd.sh'
uci set unbound.ub_main.dhcp_link='odhcpd'
uci set unbound.ub_main.dhcp4_slaac6='1'
uci set unbound.ub_main.protocol='ip6_prefer'
uci del unbound.ub_main.query_minimize
uci del unbound.ub_main.query_min_strict

uci del unbound.auth_icann
uci del unbound.fwd_isp
uci del unbound.fwd_google
uci del unbound.fwd_cloudflare

uci add unbound zone # =cfg04dc81
uci set unbound.@zone[-1].fallback='0'
uci set unbound.@zone[-1].enabled='1'
uci set unbound.@zone[-1].zone_type='forward_zone'
uci set unbound.@zone[-1].dns_assist='unprotected-loop'
uci add_list unbound.@zone[-1].server='::1'
uci add_list unbound.@zone[-1].zone_name='.'
uci set unbound.@zone[-1].port='5335'
uci set unbound.ub_main.unbound_control='1'
uci set unbound.ub_main.resource='medium'
uci set unbound.ub_main.extended_stats='1'

uci set dhcp.odhcpd.maindhcp='1'

uci set dhcp.lan.ra_dns='0'
uci set dhcp.lan.leasetime='5m'
uci set dhcp.lan.ra_mininterval='0'
uci set dhcp.lan.preferred_lifetime='5m'
uci set dhcp.lan.ra_maxinterval='120'
uci set dhcp.lan.ra_mininterval='40'
uci set dhcp.lan.ra_lifetime='360'
uci del dhcp.lan.ra_flags
uci add_list dhcp.lan.ra_flags='other-config'

uci del firewall.@defaults[0].flow_offloading='0'

uci commit

## make default_radio1 as the iot network
# /etc/config/network
uci set network.iot=interface
uci set network.iot.proto='none'

uci set network.iot.proto='static'
uci add_list network.iot.ipaddr='192.168.3.1/24'

# /etc/config/dhcp
uci set dhcp.iot=dhcp
uci set dhcp.iot.interface='iot'
uci set dhcp.iot.start='100'
uci set dhcp.iot.limit='150'
uci set dhcp.iot.leasetime='12h'

uci set dhcp.iot.dhcpv4='server'
uci set dhcp.iot.domain='iot'

# /etc/config/firewall
uci add firewall zone # =cfg10dc81
uci set firewall.@zone[-1].name='iot'
uci set firewall.@zone[-1].input='REJECT'
uci set firewall.@zone[-1].output='ACCEPT'
uci set firewall.@zone[-1].forward='REJECT'
uci add_list firewall.@zone[-1].network='iot'
uci add firewall forwarding # =cfg11ad58
uci set firewall.@forwarding[-1].src='iot'
uci set firewall.@forwarding[-1].dest='wan'

uci add firewall rule # =cfg1292bd
uci set firewall.@rule[-1].name='Iot DNS'
uci set firewall.@rule[-1].src='iot'
uci set firewall.@rule[-1].dest_port='53'
uci set firewall.@rule[-1].target='ACCEPT'
uci add firewall rule # =cfg1392bd
uci set firewall.@rule[-1].name='Iot DHCP'
uci add_list firewall.@rule[-1].proto='udp'
uci set firewall.@rule[-1].src='iot'
uci set firewall.@rule[-1].dest_port='67-68'
uci set firewall.@rule[-1].target='ACCEPT'

# /etc/config/wireless
uci set wireless.default_radio1.network='iot'

uci commit

uci del unbound.ub_main.iface_lan
uci add_list unbound.ub_main.iface_lan='iot'
uci add_list unbound.ub_main.iface_lan='lan'
uci del unbound.ub_main.iface_trig
uci add_list unbound.ub_main.iface_trig='iot'
uci add_list unbound.ub_main.iface_trig='lan'
uci add_list unbound.ub_main.iface_trig='wan'

uci commit
