# files/etc/uci-defaults/uci-settings
uci set system.@system[0].timezone=CST-8
uci set system.@system[0].zonename=Asia/Shanghai
uci set system.@system[0].hostname='link-to-the-future'
uci set system.ntp.use_dhcp='0'
uci set system.@system[0].urandom_seed='1'
uci commit system

uci set network.globals.ula_prefix='fd0f:7a7c:a8c9::/48'
uci del network.wan6
uci set firewall.@defaults[0].tcp_ecn='1'
uci set firewall.@zone[1].network='wan'
uci set network.wan.peerdns='0'
uci add_list network.wan.dns='::1'
uci set network.wan.proto='pppoe'
uci set network.wan.ipv6='auto'
uci set network.lan.igmp_snooping='1'
uci add_list network.lan.dns='fd0f:7a7c:a8c9::1'
uci add_list network.lan.dns='10.0.0.1'
uci set network.@device[0].igmp_snooping='1'
uci set network.@device[0].arp_accept='1'
uci set network.@device[0].packet_steering=0
uci set network.@device[0].flow_offloading=0
uci set network.@device[0].flow_offloading_hw=0
uci commit network

uci set dhcp.lan.ra_maxinterval='124'
uci set dhcp.lan.ra_mininterval='31'
uci set dhcp.lan.ra_lifetime='496'
uci set dhcp.lan.ra_preference='high'


# Enable wifi by default
uci set wireless.radio0.country='PA'
uci set wireless.radio1.country='PA'
uci set wireless.radio0.disabled=0
uci set wireless.radio1.disabled=0
uci set wireless.radio1.htmode='HE40'
uci set wireless.radio0.htmode='HE160'
uci set wireless.default_radio0.ssid='Link to the FUTURE'
uci set wireless.default_radio1.ssid='Xiaomi_391F'
uci set wireless.default_radio1.encryption='psk2'
uci set wireless.default_radio1.key='1qaz@WSX'
uci set wireless.default_radio0.encryption='sae'
uci set wireless.default_radio0.key='1qaz@WSX'
uci set wireless.default_radio0.ocv='0'
uci set wireless.radio1.txpower='24'
uci set wireless.radio0.txpower='21'
uci set wireless.radio0.channel='64'
uci set wireless.radio0.cell_density='0'
uci set wireless.radio1.cell_density='0'
uci set wireless.default_radio0.dtim_period='3'
uci set wireless.default_radio0.wnm_sleep_mode='1'
uci set wireless.default_radio0.multicast_to_unicast_all='1'
uci set wireless.default_radio1.dtim_period='3'
uci set wireless.default_radio1.wnm_sleep_mode='1'
uci set wireless.default_radio1.multicast_to_unicast_all='1'
uci set wireless.default_radio1.ocv='0'
uci set wireless.radio0.noscan='1'
# Force 40Mhz on 2.4Ghz
uci set wireless.radio1.noscan='1'
uci set wireless.radio1.channel='6'

uci set wireless.default_radio0.beacon_int='97'
uci set wireless.default_radio0.bss_transition='1'
uci set wireless.default_radio0.ft_over_ds='0'
uci set wireless.default_radio0.ieee80211r='1'
uci set wireless.default_radio0.ieee80211k='1'
uci set wireless.default_radio0.proxy_arp='1'
uci set wireless.default_radio0.time_advertisement='2'
uci set wireless.default_radio0.time_zone='GMT0'

uci set wireless.default_radio1.bss_transition='1'
uci set wireless.default_radio1.ft_over_ds='0'
uci set wireless.default_radio1.ieee80211r='1'
uci set wireless.default_radio1.ieee80211k='1'
uci set wireless.default_radio1.time_advertisement='2'
uci set wireless.default_radio1.time_zone='GMT0'
uci set wireless.default_radio1.max_inactivity='86400'

uci set dhcp.odhcpd.leasetrigger='/usr/lib/unbound/odhcpd.sh'
uci set unbound.ub_main.dhcp_link='odhcpd'
uci set unbound.ub_main.dhcp4_slaac6='0'
uci set unbound.ub_main.protocol='mixed'
uci del unbound.ub_main.query_minimize
uci del unbound.ub_main.query_min_strict

uci del unbound.auth_icann
uci del unbound.fwd_isp
uci del unbound.fwd_google
uci del unbound.fwd_cloudflare

uci add unbound zone # =cfg04dc81
uci set unbound.@zone[-1].fallback='0'
uci set unbound.@zone[-1].enabled='1'
uci set unbound.@zone[-1].zone_type='forward_zone'
uci set unbound.@zone[-1].dns_assist='unprotected-loop'
uci add_list unbound.@zone[-1].server='::1'
uci add_list unbound.@zone[-1].zone_name='.'
uci set unbound.@zone[-1].port='5335'
uci set unbound.ub_main.unbound_control='1'
uci set unbound.ub_main.resource='large'
uci set unbound.ub_main.recursion='aggressive'
uci set unbound.ub_main.extended_stats='1'
uci set unbound.@unbound[0].add_local_fqdn='3'
uci set unbound.@unbound[0].add_wan_fqdn='2'
uci set unbound.@unbound[0].exclude_ipv6_ga='1'

uci set dhcp.odhcpd.maindhcp='1'

uci set firewall.@defaults[0].flow_offloading='0'
uci set firewall.@defaults[0].flow_offloading_hw='0'

uci del nginx._redirect2ssl
uci del nginx._lan.listen
uci add_list nginx._lan.listen='443 ssl default_server'
uci add_list nginx._lan.listen='[::]:443 ssl default_server'
uci add_list nginx._lan.listen='80 default_server'
uci add_list nginx._lan.listen='[::]:80 default_server' 


uci commit

## make default_radio1 as the iot network
# /etc/config/network
uci set network.iot=interface

uci set network.iot.proto='static'
uci add_list network.iot.ipaddr='192.168.3.1/24'
uci set network.iot.igmp_snooping='1'
uci add_list network.iot.dns='192.168.31.1'

# /etc/config/dhcp
uci set dhcp.iot=dhcp
uci set dhcp.iot.interface='iot'
uci set dhcp.iot.start='100'
uci set dhcp.iot.limit='150'
uci set dhcp.iot.leasetime='12h'

uci set dhcp.iot.dhcpv4='server'
uci set dhcp.iot.domain='iot'

# /etc/config/firewall
uci add firewall zone # =cfg10dc81
uci set firewall.@zone[-1].name='iot'
uci set firewall.@zone[-1].input='REJECT'
uci set firewall.@zone[-1].output='ACCEPT'
uci set firewall.@zone[-1].forward='REJECT'
uci add_list firewall.@zone[-1].network='iot'
uci add firewall forwarding # =cfg11ad58
uci set firewall.@forwarding[-1].src='iot'
uci set firewall.@forwarding[-1].dest='wan'

uci add firewall rule # =cfg1292bd
uci set firewall.@rule[-1].name='Iot DNS'
uci set firewall.@rule[-1].src='iot'
uci set firewall.@rule[-1].dest_port='53'
uci set firewall.@rule[-1].target='ACCEPT'
uci add firewall rule # =cfg1392bd
uci set firewall.@rule[-1].name='Iot DHCP'
uci add_list firewall.@rule[-1].proto='udp'
uci set firewall.@rule[-1].src='iot'
uci set firewall.@rule[-1].dest_port='67-68'
uci set firewall.@rule[-1].target='ACCEPT'

uci add firewall rule # =cfg1592bd
uci set firewall.@rule[-1].name='Iot ICMP/IGMP'
uci add_list firewall.@rule[-1].proto='icmp'
uci add_list firewall.@rule[-1].proto='igmp'
uci set firewall.@rule[-1].src='iot'
uci set firewall.@rule[-1].target='ACCEPT'

uci add firewall rule # =cfg1592bd
uci set firewall.@rule[-1].name='Iot DHCPv6'
uci add_list firewall.@rule[-1].proto='udp'
uci set firewall.@rule[-1].src='iot'
uci set firewall.@rule[-1].dest_port='546-547'
uci set firewall.@rule[-1].target='ACCEPT'

uci add firewall rule # =cfg1692bd
uci set firewall.@rule[-1].name='Iot mDNS'
uci add_list firewall.@rule[-1].proto='udp'
uci set firewall.@rule[-1].src='iot'
uci set firewall.@rule[-1].target='ACCEPT'
uci set firewall.@rule[-1].dest_port='5353'

# /etc/config/wireless
uci set wireless.default_radio1.network='iot'
uci commit

uci del unbound.ub_main.iface_lan
uci add_list unbound.ub_main.iface_lan='iot'
uci add_list unbound.ub_main.iface_lan='lan'
uci del unbound.ub_main.iface_trig
uci add_list unbound.ub_main.iface_trig='iot'
uci add_list unbound.ub_main.iface_trig='lan'
uci add_list unbound.ub_main.iface_trig='wan'

uci commit
